# roles/unattended-upgrades/templates/arch/arch-auto-update.service.j2
[Unit]
Description=Automatic system updates for Arch Linux
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/arch-auto-update.sh
User=root
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/log /var/run /usr/local/bin
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes